# scripts/scheduler_examples/steward-control.service.example
[Unit]
Description=Run one steward control projection generation

[Service]
Type=oneshot
WorkingDirectory=%h/path/to/astrid
ExecStart=%h/path/to/astrid/scripts/scheduler_examples/run_projection.sh
