# scripts/com.astrid.steward-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.astrid.steward-loop</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/v/other/astrid/scripts/steward_loop_run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>7</integer></dict>
        <dict><key>Minute</key><integer>38</integer></dict>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/v/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>STEWARD_LOOP_MODEL</key>
        <string>opus</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/v/other/astrid/workspace/logs/steward_loop.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/v/other/astrid/workspace/logs/steward_loop.launchd.err</string>
    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
