# launchd/com.astrid.perception-host-ascii.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.astrid.perception-host-ascii</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/v/other/astrid/scripts/launchd_astrid_perception_host_ascii.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/v/other/astrid/capsules/perception</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/astrid-perception-host-ascii.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/astrid-perception-host-ascii.log</string>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>LOOK_SOURCE</key>
		<string>host</string>
		<key>MINIME_WORKSPACE</key>
		<string>/Users/v/other/minime/workspace</string>
		<key>ASTRID_HOST_ASCII_INTERVAL</key>
		<string>45</string>
	</dict>
</dict>
</plist>
